# deploy/alreadyapp-unsubscribed-report.timer
[Unit]
Description=Run Already Done user segment report daily at 12:00 UTC (05:00 PDT / 04:00 PST)

[Timer]
OnCalendar=*-*-* 12:00:00
Persistent=true
Unit=alreadyapp-unsubscribed-report.service

[Install]
WantedBy=timers.target
