# deploy/alreadyapp-unsubscribed-report.timer
[Unit]
Description=Run Already Done unsubscribed users report daily at 08:00 UTC

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true
Unit=alreadyapp-unsubscribed-report.service

[Install]
WantedBy=timers.target
