# deploy/alreadyapp-unsubscribed-report.service
[Unit]
Description=Already Done daily unsubscribed users report (SendGrid)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/root/alreadydone_report_service
EnvironmentFile=-/root/alreadydone_report_service/.env
ExecStart=/root/alreadydone_report_service/venv/bin/python daily_unsubscribed_report.py
StandardOutput=append:/root/alreadydone_report_service/logs/report.log
StandardError=append:/root/alreadydone_report_service/logs/report.log

[Install]
WantedBy=multi-user.target
