# deploy/alreadyapp-unsubscribed-report.timer
[Unit]
Description=Run Already Done user segment report daily at 05:00 Pacific (America/Los_Angeles)

[Timer]
OnCalendar=*-*-* 05:00:00
Timezone=America/Los_Angeles
Persistent=true
Unit=alreadyapp-unsubscribed-report.service

[Install]
WantedBy=timers.target
